# yocto/openswitch/meta-distro-openswitch/recipes-ops/utils/ops-bufmond/bufmond.service
# Copyright (C) 2015 Hewlett Packard Enterprise Development LP

[Unit]
Description=OpenSwitch Buffer monitoring Daemon
After=ovsdb-server.service ops-init.service

[Service]
Type=oneshot
CPUSchedulingPolicy=rr
CPUSchedulingPriority=10
ExecStartPre=/bin/rm -f /var/run/ops_bufmond.pid
ExecStart=/usr/bin/ops_bufmond
PIDFile=/var/run/ops_bufmond.pid

[Install]
WantedBy=multi-user.target
